# target/linux/ramips/base-files/etc/uci-defaults/01-wireless
#!/bin/sh
uci set wireless.@wifi-device[0].disabled=0
uci set wireless.@wifi-iface[0].ssid=VoCore2
uci set wireless.@wifi-iface[0].encryption='psk'
uci set wireless.@wifi-iface[0].key='12345678'
uci commit wireless
exit 0
